# scripts/crontab

*/5 * * * * PATH-TO-SCRIPTS/on_line.sh > /home/pi/log_rpi.txt 2>&1
*/3 * * * * PATH-TO-SCRIPTS/chech_flash.sh > /dev/null 2>&1
* * * * * PATH-TO-SCRIPTS/sms_check.sh > /dev/null 2>&1
